# supervisord_test.conf
[supervisord]
nodaemon=false
logfile=/home/user/webapp/supervisord_test.log
pidfile=/home/user/webapp/supervisord_test.pid

[unix_http_server]
file=/home/user/webapp/supervisor_test.sock

[supervisorctl]
serverurl=unix:///home/user/webapp/supervisor_test.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:testserver]
command=python3 /home/user/webapp/test_server.py
directory=/home/user/webapp
autostart=true
autorestart=true
stdout_logfile=/home/user/webapp/testserver.log
stderr_logfile=/home/user/webapp/testserver_error.log
